# smoke-notifier.service
[Unit]
Description=smoke-notifier — SmokePing RRD Monitor / Telegram Alert
Documentation=https://github.com/BadRush/smoke-notifer
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
User=root
Group=root

WorkingDirectory=/opt/smoke-notifier
ExecStart=/usr/bin/python3 -m smoke_notifier --config /opt/smoke-notifier/config.yaml
ExecReload=/bin/kill -HUP $MAINPID

# Restart policy
Restart=on-failure
RestartSec=10

# Environment — load .env file (- prefix = no error if missing)
EnvironmentFile=-/opt/smoke-notifier/.env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=smoke-notifier

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/smoke-notifier /tmp/smoke-notifier
PrivateTmp=false

[Install]
WantedBy=multi-user.target
